# deploy/tarkamcp.service
[Unit]
Description=TarkaMCP - Proxmox MCP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tarkamcp
EnvironmentFile=/opt/tarkamcp/.env
ExecStart=/opt/tarkamcp/.venv/bin/python -m tarkamcp serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
